# files/restic/backup.timer
[Unit]
Description=Make a backup with restic daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
