# files/freedns/freedns.timer
[Timer]
OnCalendar=*-*-* *:00,05,10,15,20,25,30,35,40,45,50,55:47
Persistent=true

[Install]
WantedBy=timers.target
