# files/deduplication/dedup.timer
[Unit]
Description=Remove duplicated files monthly

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
